# systemd/seclyzer-mouse@.service
[Unit]
Description=SecLyzer Mouse Collector
After=redis.service
Wants=redis.service

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/Documents/Projects/SecLyzer
ExecStart=/home/%i/Documents/Projects/SecLyzer/collectors/mouse_collector/target/release/mouse_collector
Restart=always
RestartSec=10
StandardOutput=append:/var/log/seclyzer/mouse_collector.log
StandardError=append:/var/log/seclyzer/mouse_collector.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
